# scripts/com.lo-first.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.lo-first</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/lo-first</string>
        <string>-config</string>
        <string>/etc/lo-first/config.yaml</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/var/log/lo-first.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/lo-first.error.log</string>
</dict>
</plist>